# ATUALIZACAO_CALCULO-main/ATUALIZACAO_CALCULO-main/deploy/atualizacao-calculo-api.service
[Unit]
Description=API interna de actualização de cálculo (precatórios / planilha / memoria_calculo)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Ajuste os caminhos abaixo se o repositório não estiver neste volume.
User=www-data
Group=www-data
WorkingDirectory=/mnt/volume_nyc1_1778499775066/Centro_Organizacional/ATUALIZACAO_CALCULO-main/ATUALIZACAO_CALCULO-main
EnvironmentFile=-/mnt/volume_nyc1_1778499775066/Centro_Organizacional/ATUALIZACAO_CALCULO-main/ATUALIZACAO_CALCULO-main/.env
ExecStart=/mnt/volume_nyc1_1778499775066/Centro_Organizacional/.venv/bin/python /mnt/volume_nyc1_1778499775066/Centro_Organizacional/ATUALIZACAO_CALCULO-main/ATUALIZACAO_CALCULO-main/api_atualizacao_calculo.py
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
